# packages/simplyprint-duet3d/simplyprint_duet3d-1.3.11-py3-none-any.whl/simplyprint-duet3d.service
[Unit]
Description=Meltingplot Duet <-> SimplyPrint.io connector
After=network.target

[Service]
Type=simple
Restart=always
StartLimitInterval=0
StartLimitBurst=1440
RestartSec=60
ExecStart=/usr/local/bin/simplyprint-duet3d start
WorkingDirectory=/usr/local
User=ubuntu
Group=ubuntu

[Install]
WantedBy=multi-user.target