# packages/gopnik/gopnik-0.1.0.tar.gz/gopnik-0.1.0/docker/api/supervisord.conf
# Supervisor configuration for Gopnik API server

[supervisord]
nodaemon=true
user=root
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
user=root
autostart=true
autorestart=true
stdout_logfile=/app/logs/nginx_stdout.log
stderr_logfile=/app/logs/nginx_stderr.log
priority=10

[program:gopnik-api]
command=gunicorn -c gunicorn.conf.py src.gopnik.interfaces.api.app:app
directory=/app
user=gopnik
autostart=true
autorestart=true
stdout_logfile=/app/logs/api_stdout.log
stderr_logfile=/app/logs/api_stderr.log
priority=20
environment=GOPNIK_ENV=production,GOPNIK_CONFIG=/app/config/gopnik.yaml

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface